# atum-desk/infra/systemd/atum-desk-api.service
[Unit]
Description=ATUM DESK API Service - Production Hardened
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=notify
User=navi
Group=navi
WorkingDirectory=/data/ATUM DESK/atum-desk/api
Environment="PATH=/data/ATUM DESK/atum-desk/api/.venv/bin"
Environment="PYTHONPATH=/data/ATUM DESK/atum-desk/api"
Environment="ENVIRONMENT=production"
Environment="DATABASE_URL=postgresql+asyncpg://postgres:postgres@localhost:5432/atum_desk"
# Redis removed per NO REDIS constraint
ExecStart=/data/ATUM DESK/atum-desk/api/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000 --workers 2
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
Restart=on-failure
RestartSec=5
StartLimitBurst=5
StartLimitIntervalSec=60

# Resource limits
MemoryMax=2G
MemoryHigh=1G
CPUQuota=150%

# File descriptor limit
LimitNOFILE=65536

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ReadOnlyPaths=/data/ATUM DESK/atum-desk/api

[Install]
WantedBy=multi-user.target
