# atum-desk/infra/systemd/atum-desk-api.service
[Unit]
Description=ATUM DESK API Service
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service

[Service]
Type=simple
User=navi
Group=navi
WorkingDirectory=/data/ATUM DESK/atum-desk/api
Environment="PATH=/data/ATUM DESK/atum-desk/api/.venv/bin"
Environment="PYTHONPATH=/data/ATUM DESK/atum-desk/api"
Environment="ENVIRONMENT=production"
Environment="DATABASE_URL=postgresql+asyncpg://postgres:postgres@localhost:5432/atum_desk"
Environment="REDIS_URL=redis://localhost:6379/0"
Environment="SECRET_KEY=atum-desk-production-secret-key-change-this"
ExecStart=/data/ATUM DESK/atum-desk/api/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000 --workers 2
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
